# services/backup.service
[Unit]
Description=Backup db

[Service]
Type=oneshot
ExecStart=/bin/bash /home/debian/cycling-routes/scripts/backup.sh